# uav.service
[Unit]
Description=Uav Service
After=network.target
 
[Service]
User=user
ExecStart=/javodata/run.sh
 
[Install]
WantedBy=multi-user.target